# profile.d/local05-setupstack.sh
#!/bin/sh
if [ -n "${RUNNING_INSIDE_JUPYTERLAB}" ]; then
    . /opt/romancal/runtime/loadromancal
fi
